# pump-bot.service
[Unit]
Description=MEXC Pump & Dump Scanner Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Замени YOUR_USER на своего пользователя сервера (например: ubuntu, root, deploy)
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/pump-bot
EnvironmentFile=/home/YOUR_USER/pump-bot/.env
ExecStart=/home/YOUR_USER/pump-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
